# rc.d/bhyved.in
#!/bin/sh

# PROVIDE: bmd
# REQUIRE: DAEMON syslogd
# BEFORE: LOGIN
# KEYWORD: nojail shutdown

. /etc/rc.subr

name="bmd"
desc="bhyve management daemon"
rcvar="bmd_enable"
command="%%LOCALBASE%%/sbin/${name}"
pidfile="/var/run/bmd.pid"
extra_commands="reload"
required_modules="vmm"

load_rc_config $name
run_rc_command "$1"
